# docker/etc/s6-overlay/s6-rc.d/init-data-permissions/up
/etc/s6-overlay/scripts/init-data-permissions.sh