# docker/etc/s6-overlay/s6-rc.d/grimoire/run
#!/command/with-contenv bash

cd /app

# Run migrations first
s6-setuidgid grimoire /usr/local/bin/bun run run-migrations

# Start the application
exec s6-setuidgid grimoire /usr/local/bin/bun ./build/index.js
